# bilibili-sync-api.service
[Unit]
Description=Bilibili Sync API Service
Documentation=https://github.com/your-repo/bilibili-sync-app
After=network.target

[Service]
# ========== 服务类型 ==========
Type=notify
User=www-data
Group=www-data
WorkingDirectory=/var/www/bilibili-sync-app

# ========== 环境配置 ==========
Environment="PATH=/var/www/bilibili-sync-app/venv/bin"
EnvironmentFile=/var/www/bilibili-sync-app/.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bilibili-sync

# ========== 进程管理 ==========
ExecStart=/var/www/bilibili-sync-app/venv/bin/gunicorn \
    --workers 4 \
    --worker-class uvicorn.workers.UvicornWorker \
    --worker-tmp-dir /dev/shm \
    --bind 127.0.0.1:8000 \
    --bind unix:/var/run/bilibili-sync-api.sock \
    --access-logfile /var/www/bilibili-sync-app/logs/access.log \
    --error-logfile /var/www/bilibili-sync-app/logs/error.log \
    --log-level info \
    --timeout 300 \
    --keep-alive 5 \
    --max-requests 1000 \
    --max-requests-jitter 100 \
    main:app

ExecReload=/bin/kill -s HUP $MAINPID

# ========== 重启策略 ==========
Restart=always
RestartSec=10

# ========== 启动延迟 ==========
StartLimitInterval=600
StartLimitBurst=3

# ========== 资源限制 ==========
# 最大打开文件数
LimitNOFILE=65535

# 最大进程数
LimitNPROC=65535

# CPU 时间限制 (无限制)
CPUQuota=100%

# 内存限制 (可根据需要调整)
MemoryMax=2G

# ========== 安全配置 ==========
# 不允许新特权
NoNewPrivileges=true

# 私有临时目录
PrivateTmp=true

# 保护系统文件
ProtectSystem=strict

# 保护主目录
ProtectHome=true

# 设备访问保护
PrivateDevices=true

# IPC 隔离
PrivateIPC=true

# 只读路径
ReadOnlyPaths=/etc /usr /lib /lib64 /sys

# 可写路径
ReadWritePaths=/var/www/bilibili-sync-app

# ========== 日志和追踪 ==========
# 启用系统日志
SyslogIdentifier=bilibili-sync-api

# ========== 特权下降 ==========
# 确保 www-data 用户运行此服务
DynamicUser=false

[Install]
WantedBy=multi-user.target
