# bilibili-sync-api.service
[Unit]
Description=Bilibili Sync API Service
Documentation=https://github.com/your-repo/bilibili-sync-app
After=network.target

[Service]
# ========== 服务类型 ==========
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/bilibili-sync-app

# ========== 环境配置 ==========
EnvironmentFile=/opt/bilibili-sync-app/.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bilibili-sync

# ========== 进程管理 ==========
ExecStart=/opt/bilibili-sync-app/bilibili-sync-api

# ========== 重启策略 ==========
Restart=always
RestartSec=10

# ========== 启动延迟 ==========
StartLimitInterval=600
StartLimitBurst=3

# ========== 资源限制 ==========
# 最大打开文件数
LimitNOFILE=65535

# 最大进程数
LimitNPROC=65535

# CPU 时间限制 (无限制)
CPUQuota=100%

# 内存限制 (可根据需要调整)
MemoryMax=2G

# ========== 安全配置 ==========
# 不允许新特权
NoNewPrivileges=true

# 私有临时目录
PrivateTmp=true

# 保护系统文件
ProtectSystem=strict

# 保护主目录
ProtectHome=true

# 设备访问保护
PrivateDevices=true

# IPC 隔离
PrivateIPC=true

# 只读路径
ReadOnlyPaths=/etc /usr /lib /lib64 /sys

# 可写路径
ReadWritePaths=/var/www/bilibili-sync-app

# ========== 日志和追踪 ==========
# 启用系统日志
SyslogIdentifier=bilibili-sync-api

# ========== 特权下降 ==========
# 确保 www-data 用户运行此服务
DynamicUser=false

[Install]
WantedBy=multi-user.target
